# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:knowledge-base]
command=/path/to/venv/bin/gunicorn -c /path/to/knowledge-base/gunicorn_config.py app:app
directory=/path/to/knowledge-base
user=your-username
autostart=true
autorestart=true
stderr_logfile=/var/log/knowledge-base/err.log
stdout_logfile=/var/log/knowledge-base/out.log
environment=FLASK_ENV="production",DATABASE_URL="postgresql://user:password@localhost:5432/knowledge_base"